# profiles/user-path/root/etc/profile.d/local-bin.sh
case ":$PATH:" in
  *":$HOME/.local/bin:"*) ;;
  *) PATH="$PATH:$HOME/.local/bin" ;;
esac
